# osg-info-services.cron
# Cron job
# The lock file can be enabled or disabled via a
# service osg-info-services start
# chkconfig osg-info-services on

http_proxy=""

# Run every 10 minutes +/- 5 minutes (300 sec)
# Send soft kill after 3 minutes (180 sec)
# Send hard kill 1.5 minutes after that (90 sec)
*/10 * * * * root /usr/libexec/osg-info-services/cronjob-wrapper /var/lock/subsys/osg-info-services 300 /usr/libexec/osg-info-services/run-with-timeout 180:90 /usr/sbin/osg-info-services
